# systemd/ralph-stock.service
# Ralph Stock Analysis - systemd Service Unit
# =============================================================================
# 재부팅 시 Docker Compose 서비스 자동 시작
# =============================================================================
#
# 설치 방법:
#   sudo cp systemd/ralph-stock.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable ralph-stock.service
#   sudo systemctl start ralph-stock.service
#
# 관리 명령어:
#   sudo systemctl status ralph-stock    # 상태 확인
#   sudo systemctl restart ralph-stock   # 재시작
#   sudo systemctl stop ralph-stock      # 중지
#   sudo journalctl -u ralph-stock -f    # 로그 보기
# =============================================================================

[Unit]
Description=Ralph Stock Analysis - Docker Compose Services (Production)
Documentation=https://github.com/ralphpark/kr_stock_analysis
After=docker.service network-online.target
Wants=docker.service network-online.target
Requires=docker.service

# 서비스 시작 전 네트워크가 준비될 때까지 대기
# docker-compose가 네트워크에 연결된 컨테이너를 시작할 수 있도록 함

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/ralph/work/python/kr_stock_analysis
ExecStartPre=/bin/bash -c 'docker network inspect ralph-network >/dev/null 2>&1 || docker network create ralph-network'
ExecStart=/usr/bin/docker compose --profile prod up -d
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300
TimeoutStopSec=60

# 환경 변수
Environment="COMPOSE_PROFILES=prod"

# 재시작 정책
Restart=on-failure
RestartSec=10

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralph-stock

[Install]
WantedBy=multi-user.target
